# packages/gateway/deploy/chatcode-gateway.service
[Unit]
Description=Chatcode.dev Gateway Daemon
Documentation=https://chatcode.dev/docs/gateway
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
User=vibe
Group=vibe
WorkingDirectory=/home/vibe
ExecStart=/usr/local/bin/chatcode-gateway
EnvironmentFile=/etc/chatcode/gateway.env
Restart=on-failure
RestartSec=5s

# Keep service environment close to host defaults so explicitly-approved sudo
# commands behave the same as on a normal SSH shell.
PrivateTmp=false

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=chatcode-gateway

[Install]
WantedBy=multi-user.target
